# rpi_servers/rpi_mirror/etc/cron.d/usb-mirror-raspbian
20 4	* * *	root	/usr/local/sbin/usb-mirror-raspbian
